# 99-udrivemount.rules
ENV{ID_PATH}=="*usb*", ENV{DEVTYPE}=="partition", RUN+="/bin/systemctl restart udrivemount@%k.service"
