# 99-udrivemount.rules
ENV{ID_USB_DRIVER}=="usb-storage", ENV{DEVTYPE}=="partition", RUN+="/bin/systemctl restart udrivemount@%k.service"
